# release_hi3521/host/rootfs/etc/udev/rules.d/91-usb_mouse.rules
ACTION!="add", KERNELS=="input*", SUBSYSTEMS=="input", RUN+="/usr/sbin/usbhid_reload.sh"
